# payload/skel/.config/systemd/user/my-desk-wallpaper-random.timer
[Unit]
Description=Random my-desk wallpaper rotation

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
RandomizedDelaySec=2min
Persistent=true
Unit=my-desk-wallpaper-random.service

[Install]
WantedBy=timers.target
